# scheduled-jobs/cert-monitor/bloom-cert-monitor.timer
[Unit]
Description=Run bloom-cert-monitor weekly (Sunday 09:00 server time)
Documentation=https://github.com/Salk-Harnessing-Plants-Initiative/bloom/issues/285

[Timer]
OnCalendar=Sun 09:00
Persistent=true
Unit=bloom-cert-monitor.service

[Install]
WantedBy=timers.target
